# sys-process/glances/files/glances-web.service
[Unit]
Description=Glances Web Server

[Service]
ExecStart=/usr/bin/glances -w

[Install]
WantedBy=multi-user.target
